# timetag-acquire.rules
SUBSYSTEMS=="usb", ACTION=="add", ATTR{idVendor}=="04b4", ATTR{idProduct}=="1004", RUN+="/usr/bin/timetag_acquire -s /tmp/timetag.sock -d"
